# deploy/serverpilot-admin.service
[Unit]
Description=ServerPilot Admin Frontend (ServerHQ)
After=network.target
Wants=serverpilot-server-hq.service

[Service]
Type=simple
User=serverpilot
Group=serverpilot
WorkingDirectory=/opt/serverpilot
Environment=NODE_ENV=production
ExecStart=/usr/bin/npx next start -p 3000
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=serverpilot-admin

[Install]
WantedBy=multi-user.target
